# faster.service
[Unit]
Description=Faster
After=multi-user.target

[Service]
ExecStart=/usr/bin/faster

[Install]
WantedBy=default.target
